# ops/watchdog/smarthome-watchdog.service
[Unit]
Description=Smart Home System Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /home/trungcoolboy/smart-home/ops/watchdog/system_watchdog.py /var/log/smarthome-watchdog
Restart=always
RestartSec=2
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
